# ubuntu/src/dockerstartup/modes/sshd.conf
[program:ssh]
command=/usr/sbin/sshd -D
stdout_events_enabled=true
stderr_events_enabled=true
autorestart=true