# ubuntu/src/modes/vncserver.conf
[program:vncserver]
command=vncserver -geometry 1152x864 :1 securitytypes=none
stdout_events_enabled=true
stderr_events_enabled=true
autorestart=true